# scripts/rlight.service
[Unit]
Description=RLight Portaria Digital v7 - Backend
After=network.target

[Service]
Type=simple
User=portaria
WorkingDirectory=/home/portaria/rlight/host
ExecStart=/home/portaria/rlight/venv/bin/python main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rlight-backend

[Install]
WantedBy=multi-user.target
