# scripts/rlight-ui.service
[Unit]
Description=RLight Portaria Digital v7 - UI Kiosk
After=rlight.service graphical.target
Wants=graphical.target

[Service]
Type=simple
User=portaria
WorkingDirectory=/home/portaria/rlight
Environment=DISPLAY=:0
ExecStart=/bin/bash /home/portaria/rlight/scripts/v7_kiosk.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rlight-ui

[Install]
WantedBy=graphical.target
